# bash/systemd-template/fastapi-conda.service
[Unit]
Description=FastAPI app (Conda env)
After=network.target

[Service]
User=youruser
Group=youruser
WorkingDirectory=/path/to/fastapi
ExecStart=/bin/bash -c 'source /home/youruser/miniconda3/etc/profile.d/conda.sh && conda activate yourenv && exec uvicorn main:app --host 0.0.0.0 --port 8000'
Restart=always

[Install]
WantedBy=multi-user.target
